# dotfiles/systemd/user/kwallet.service
[Unit]
Description=KWallet Daemon
After=graphical-session.target

[Service]
ExecStart=/usr/bin/kwalletd5

[Install]
WantedBy=default.target
